# deploy/helloasso.service
[Unit]
Description=HelloAsso Manager
After=network.target

[Service]
User=<your-user>
Group=<your-user>
WorkingDirectory=/home/<your-user>/helloAssoManager
EnvironmentFile=/etc/helloasso/env
ExecStart=/home/<your-user>/helloAssoManager/.venv/bin/gunicorn \
    --workers 3 \
    --bind unix:/run/helloasso/gunicorn.sock \
    config.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
RuntimeDirectory=helloasso
Restart=on-failure

[Install]
WantedBy=multi-user.target
